# jetson/agrikd.service
[Unit]
Description=AgriKD Edge Inference Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=agrikd
WorkingDirectory=/opt/agrikd
ExecStart=/usr/bin/python3 /opt/agrikd/app/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

# Resource limits
MemoryMax=512M
CPUQuota=80%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/agrikd

[Install]
WantedBy=multi-user.target
