# jetson/agrikd.service
# TEMPLATE — setup_jetson.sh generates the actual service file.
# This file is a reference for manual installs (copy and adjust paths).

[Unit]
Description=AgriKD Edge Inference Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agrikd
WorkingDirectory=/opt/agrikd
ExecStart=/usr/bin/python3 /opt/agrikd/app/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agrikd
Environment=PYTHONUNBUFFERED=1

# Crash loop protection: max 5 restarts in 300 seconds
StartLimitBurst=5
StartLimitIntervalSec=300

# Resource limits
MemoryMax=512M
CPUQuota=80%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
ReadWritePaths=/opt/agrikd

[Install]
WantedBy=multi-user.target
